# gbmc-s5-poweroff.target
[Unit]
Description=Sys Host S5 Power Off
