# gbmc-psu-hardreset-pre.target
[Unit]
Description=Sys PSU Hard Reset
Wants=gbmc-psu-hardreset-time.service
After=gbmc-psu-hardreset-time.service
